# deploy/pyspire-python.service
# PySpire - static site served by Python's built-in http.server on port 8081.
#
# PySpire is a pure static site (index.html + css/ + js/), so all this unit
# does is serve those files. Zero dependencies beyond python3, which ships
# with virtually every Linux distro.
#
# Install:
#   sudo mkdir -p /var/www/pyspire
#   sudo cp -r index.html css js /var/www/pyspire/
#   sudo cp deploy/pyspire.service /etc/systemd/system/pyspire.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pyspire.service
#
# Verify:
#   systemctl status pyspire.service
#   curl -I http://localhost:8081/        (expect HTTP/1.0 200 OK)
#   journalctl -u pyspire.service -f      (live access log)
#
# Confirm your python3 path with `command -v python3` and adjust ExecStart
# below if it is not /usr/bin/python3.

[Unit]
Description=PySpire static site (Python http.server) on port 8081
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# systemd allocates a throwaway unprivileged user just for this service.
# Nothing to create or manage, and it cannot log in.
DynamicUser=true

# Serve the deployed copy. Binds all interfaces so the site is reachable on
# the LAN; if you put a reverse proxy in front, change --bind to 127.0.0.1.
WorkingDirectory=/var/www/pyspire
ExecStart=/usr/bin/python3 -m http.server 8081 --bind 0.0.0.0 --directory /var/www/pyspire

# Keep it alive across crashes / OOM kills.
Restart=on-failure
RestartSec=5s

# --- Sandboxing: contain the process if it is ever compromised ---
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
# The site directory is only ever read, never written, by this service.
ReadOnlyPaths=/var/www/pyspire
# Port 8081 is unprivileged (>1024), so no extra capabilities are needed.

[Install]
WantedBy=multi-user.target
